# openuem-server_0.1.0_amd64/etc/systemd/system/openuem-nats.service
[Unit]
Description=OpenUEM NATS Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/openuem-server/bin/openuem-nats-service
User=openuem
Group=openuem

[Install]
WantedBy=network.target