# openuem-server_0.1.1_amd64/etc/systemd/system/openuem-ocsp-responder.service
[Unit]
Description=OpenUEM OCSP Responder
After=network.target

[Service]
Type=simple
ExecStart=/opt/openuem-server/bin/openuem-ocsp-responder
User=openuem
Group=openuem

[Install]
WantedBy=network.target