# openuem-server-0.1.0-amd64/etc/systemd/system/openuem-agent-worker.service
[Unit]
Description=OpenUEM Agent Worker
After=network.target

[Service]
Type=simple
ExecStart=/opt/openuem-server/bin/openuem-agent-worker
User=openuem
Group=openuem

[Install]
WantedBy=network.target